# config_portal.service
[Unit]
Description=SDR Config Portal — unified web configuration UI (port 8083)
After=network.target

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/config_portal.py 8083
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
