# noaa_capture.service
[Unit]
Description=NOAA APT Satellite Capture
After=network.target

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/noaa_capture.py
Restart=on-failure
RestartSec=30
User=orangepi

[Install]
WantedBy=multi-user.target
