# rtl_433.service
[Unit]
Description=RTL-433 — 433MHz sensor decoder
After=network.target

[Service]
ExecStartPre=/bin/mkdir -p /var/log/rtl_433
ExecStart=/usr/bin/rtl_433 \
    -f 433.92M \
    -F json:/var/log/rtl_433/events.json \
    -F http:0.0.0.0:8433
Restart=always
RestartSec=5
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
